# docker-supervisord.conf
[program:djangopackages]
command=docker-compose up
directory=/code/djangopackages --abort-on-container-exit
redirect_stderr=true
autostart=true
autorestart=true
priority=10
