# docker-supervisord.conf
[program:djangopackages]
command=docker-compose up
directory=/code/djangopackages
redirect_stderr=true
autostart=true
autorestart=true
priority=10
